# vm-systemd/systemd-random-seed.service.d/30_qubes.conf
[Unit]
After=qubes-db.service

[Service]
ExecStartPost=/usr/lib/qubes/init/qubes-random-seed.sh
